# scm/etc/scm_backup.cron
#! /bin/bash
# backup svn & trac daily

export HOME=/root
export MAILTO=maqm@ihep.ac.cn 
export ENV_HOME=$HOME/env
env-(){ . $HOME/env/env.bash && env-env $* ; }

env-
scm-backup-
type scm-backup-all

scm-backup-all
scm-backup-rsync
case $NODE_TAG in
  Y1)
    ;;
   *)
    scm-backup-monitor-ihep
    ;;
esac
scm-backup-mail

exit 0
